# root/etc/services.d/squid-logs/run
#!/usr/bin/with-contenv bash
while true; do
  if [[ -f /usr/local/squid/var/logs/access.log ]]; then
    tail -f /usr/local/squid/var/logs/access.log
  else
    sleep 1
  fi
done
