# root/etc/services.d/squid/run
#!/usr/bin/execlineb

s6-envuidgid squid
s6-applyuidgid -U

/usr/local/squid/sbin/squid --foreground -d 1 -Y 2>&1
