# root/etc/services.d/squid/run
#!/usr/bin/with-contenv bash

exec squid --foreground -d 5 -Y 2>&1
